# udev/92-charging.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/home/jerome/.config/udev/scripts/92-charging.sh disconnected"
